# config/systemd/ipu7-camera-dynamic.service.d/override.conf
[Unit]
After=local-fs.target pipewire.service
Before=display-manager.service

[Service]
ExecStartPre=+/usr/bin/modprobe intel_ipu7
ExecStartPre=+/usr/bin/modprobe intel_ipu7_isys
ExecStartPre=+/usr/bin/modprobe intel_ipu7_psys
ExecStartPre=+/usr/bin/modprobe imx471
ExecStartPre=-+/usr/bin/modprobe v4l2loopback
ExecStartPre=+/usr/bin/modprobe v4l2loopback
ExecStartPre=+/usr/bin/systemd-tmpfiles --create /etc/tmpfiles.d/camera-run.conf
ExecStartPre=+/bin/sh -c 'ipcrm -M 0x43414D 2>/dev/null || true; rm -f /dev/shm/sem.camlock /run/ipu7-camera-dynamic.lock 2>/dev/null || true'
Environment=CAMERA_OUT_FORMAT=yuy2
Environment=CAMERA_SENSOR_VFLIP=0
Environment=CAMERA_GST_VFLIP=vertical-flip
Environment=CAMERA_CHROMA=nv12
Environment=CAMERA_COLOR_FIX=none
Environment=CAMERA_AWB_MODE=auto
Environment=CAMERA_ACTIVATION_BACKOFF=10.0
Environment=CAMERA_STOP_DELAY=2.0
Environment=CAMERA_BALANCE_HUE=-0.35
Environment=CAMERA_BALANCE_SAT=1.25
Environment=CAMERA_CCM_RED=1.0294,0,0,0
Environment=CAMERA_CCM_GREEN=0,0.9621,0,0
Environment=CAMERA_CCM_BLUE=0,0,1.2892,0
Environment=CAMERA_RESOLUTION=auto
Environment=CAMERA_ALWAYS_ON=0
Environment=CAMERA_POLL_MS=2000
Environment=CAMERA_ACTIVATION_DELAY=0.8
Environment=CAMERA_PIPEWIRE=0
Environment=CAMERA_PIPEWIRE_UID=1000
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/1000
Environment=LIBCAMERA_IPA_MODULE_PATH=/usr/lib/libcamera/ipa
Environment=GST_DEBUG=
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5
